# DEBIAN/DEBIAN/usr/local/bin/etc/systemd/system/webblock.service
[Unit]
Description=WebBlock - aggiorna blocco siti web
After=network.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/webblock
